# meta-agl-container-demo/meta-container-host/recipes-containers/lxc-demo/files/lxc-watchdog-ic.service
[Unit]
Description=Trigger to reboot IC container

[Service]
Type=simple
ExecStart=/usr/bin/watchdog_ic.sh
