# meta-agl-container-demo/meta-container-host/recipes-containers/lxc-demo/files/lxc-watchdog-system.service
[Unit]
Description=Trigger to reboot the system

[Service]
Type=simple
ExecStart=/usr/bin/watchdog_system.sh
