# meta-agl-container-demo/meta-container-host/recipes-containers/lxc-demo/files/lxc-watchdog-ivi.service
[Unit]
Description=Trigger to reboot IVI container

[Service]
Type=simple
ExecStart=/usr/bin/watchdog_ivi.sh
